# supervisord-gae.ini
[program:gae]
command=bash /root/gae-env/run.sh
autostart=true
autorestart=true
stdout_logfile=/var/log/supervisord/gae.log
stderr_logfile=/var/log/supervisord/gae_error.log
